# deploy/pager.service
[Unit]
Description=pager — PagerMon log to PDF + printing
After=network.target

[Service]
Type=simple
User=pager
Group=pager
# pager must be in the `lp` group to talk to CUPS for printing.
SupplementaryGroups=lp
WorkingDirectory=/opt/pager
Environment=PAGER_CONFIG=/opt/pager/config
ExecStart=/opt/pager/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
